# cronk.service.example
[Unit]
Description=cronk
ConditionPathExists=/path/to/cronk.yaml

[Service]
ExecStart=/path/to/cronk.py

[Install]
WantedBy=default.target
